# roles/sumaris-pod/templates/systemd/sumaris-db-backup.service.j2
[Unit]
Description=Let's Encrypt {{ sumaris_host }} renewal

[Service]
Type=oneshot
ExecStart=/opt/{{ username }}/{{ target_dir }}/bin/sumaris-db.sh backup


