# roles/sumaris-pod/templates/systemd/sumaris-db-backup.timer.j2
[Unit]
Description=Daily backup of Sumaris DB {{ username }}

[Timer]
OnCalendar=daily UTC
RandomizedDelaySec=10
Persistent=true

[Install]
WantedBy=timers.target